# etc/supervisor/conf.d/elasticsearch.conf
[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch -Enetwork.host=127.0.0.1 -Ehttp.port=9200 -Etransport.tcp.port=9300
numprocs=1
autostart=true
autorestart=true
user=elasticsearch
echo environment=ES_HEAP_SIZE="20g"
environment=ES_PATH_CONF=/etc/elasticsearch