# etc/supervisor/conf.d/kibana.conf
[program:kibana]
command=/usr/share/kibana/bin/kibana -c /etc/kibana/kibana.yml
numprocs=1
autostart=true
autorestart=true
user=kibana
directory=/usr/share/kibana/
