# etc/supervisor/conf.d/kibana.conf
[program:kibana]
command=/usr/bin/kibana
numprocs=1
autostart=true
autorestart=true
user=kibana
