# etc/supervisor/conf.d/zookeeper.conf
[program:zookeeper]
command=/usr/share/zookeeper/bin/zkServer.sh start-foreground
autorestart=true
stopsignal=KILL
