# install/storm/storm-supervisor.conf
[program:storm-supervisor]
command=storm supervisor
directory=/home/storm
autorestart=true
user=storm
