# install/storm/storm-nimbus.conf
[program:storm-nimbus]
command=storm nimbus
directory=/home/storm
autorestart=true
user=storm
