# install/storm/storm-ui.conf
[program:storm-ui]
command=storm ui
directory=/home/storm
autorestart=true
user=storm
